# systemd/hydra-credential-broker.service
[Unit]
Description=Hydra Credential Broker
After=network.target

[Service]
Type=simple
User=aisvc
Group=aisvc
ExecStart=/usr/bin/python3 /opt/claude-swarm/src/credential_broker.py
Restart=on-failure
RestartSec=5s
RuntimeDirectory=hydra
RuntimeDirectoryMode=0750
Environment="SSH_AUTH_SOCK=/run/user/%U/keyring/ssh"

[Install]
WantedBy=default.target
