# deploy/swarm-dashboard.service
[Unit]
Description=claude-swarm Web Dashboard (port 9192)
After=network.target
Documentation=https://github.com/scoobydont-666/claude-swarm

[Service]
Type=simple
User=josh
Group=josh
WorkingDirectory=/opt/claude-swarm
Environment=PYTHONPATH=/opt/claude-swarm/src
ExecStart=/usr/bin/python3 /opt/claude-swarm/src/dashboard.py --host 127.0.0.1 --port 9192
Restart=on-failure
RestartSec=10

# Security hardening
ProtectSystem=strict
ReadWritePaths=/opt/swarm /opt/claude-swarm/data
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
