# deploy/swarm-health-monitor.service
[Unit]
Description=claude-swarm Health Monitor — fleet health daemon with auto-remediation
After=network.target

[Service]
Type=simple
User=josh
WorkingDirectory=/opt/claude-swarm
Environment=PYTHONPATH=/opt/claude-swarm/src
ExecStart=/usr/bin/python3 /opt/claude-swarm/src/health_monitor.py
Restart=on-failure
RestartSec=30

# Hardening
ProtectSystem=strict
ReadWritePaths=/opt/swarm /opt/claude-swarm/data
PrivateTmp=true
NoNewPrivileges=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=swarm-health-monitor

[Install]
WantedBy=multi-user.target
